# crontab
# Check for airport data updates daily at 2:00 AM
0 2 * * * cd /home/jennertorrence/code/air_field_info_server && poetry run python manage.py import_airports >> /var/log/airfield_info/airport_updates.log 2>&1

